# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:jupyter]
command=/usr/local/bin/start.sh start-notebook.py

[program:aws]
command=/usr/local/bin/aws.sh
